# deploy/systemd/detect.service
[Unit]
Description=Edge CCTV detector stack
After=ingest.service
Requires=docker.service

[Service]
Restart=always
RestartSec=5
Environment=PROFILE=all
ExecStart=/usr/bin/env bash -lc 'cd {{INSTALL_PREFIX}}/edge-cctv-jetson && PROFILE=${PROFILE} docker compose --profile all up detect_object detect_face detect_action tracker rules events privacy recorder exporter'
ExecStop=/usr/bin/env bash -lc 'cd {{INSTALL_PREFIX}}/edge-cctv-jetson && docker compose stop detect_object detect_face detect_action tracker rules events privacy recorder exporter'

[Install]
WantedBy=edge-cctv.target
