# deploy/systemd/edge-cctv.target
[Unit]
Description=Edge CCTV multi-service target
Requires=ingest.service detect.service ui.service
After=docker.service
AllowIsolate=yes
